# rootfs/usr/lib/systemd/system/fedora-nix-graphics-rebuild.service
[Unit]
Description=Rebuild system-level Nix graphics runtime when Nix config changes
After=network-online.target nix.mount nix-daemon.socket
Wants=network-online.target nix-daemon.socket
Requires=nix.mount
ConditionPathExists=/etc/nix/flake.nix

[Service]
Type=oneshot
ExecStart=/usr/libexec/fedora-nix-graphics-rebuild
RemainAfterExit=yes
