# rootfs/usr/lib/systemd/system/non-nixos-gpu-refresh.path
[Unit]
Description=Refresh /run/opengl-driver when Nix graphics driver profile changes
After=local-fs.target nix.mount
Requires=nix.mount
ConditionPathExists=/var/nix-system/graphics-drivers
ConditionPathIsDirectory=/var/nix-system/graphics-drivers

[Path]
PathChanged=/var/nix-system/graphics-drivers
Unit=non-nixos-gpu.service

[Install]
WantedBy=multi-user.target
